# config/supervisor/conf.d/g4he.conf
[program:g4he]
command=/opt/g4he/bin/python /opt/g4he/src/g4he/portality/app.py
user=cloo
directory=/opt/g4he/src/g4he
stdout_logfile=/var/log/supervisor/%(program_name)s-access.log
stderr_logfile=/var/log/supervisor/%(program_name)s-error.log
autostart=true
autorestart=true
stopasgroup=true
