# config/supervisor/conf.d/oacwellcome-production.conf
[program:oacwellcome-production]
command=/opt/oacwellcome/bin/newrelic-admin run-program /opt/oacwellcome/bin/gunicorn -c /opt/oacwellcome/src/oacwellcome/deploy/oacwellcome_gunicorn_config.py service.web:app
environment = NEW_RELIC_CONFIG_FILE=/opt/oacwellcome/src/oacwellcome/deploy/newrelic.ini,NEW_RELIC_ENVIRONMENT=production
user=cloo
directory=/opt/oacwellcome/src/oacwellcome
stdout_logfile=/var/log/supervisor/%(program_name)s-access.log
stderr_logfile=/var/log/supervisor/%(program_name)s-error.log
autostart=true
autorestart=true
stopasgroup=true
