# config/supervisor/conf.d/oagr-test-daemon.conf
[program:oagr-test-daemon]
command=/opt/lantern/bin/python /opt/lantern/src/lantern/magnificent-octopus/octopus/modules/oag/runner.py
user=cloo
directory=/opt/lantern/src/lantern
stdout_logfile=/var/log/supervisor/%(program_name)s-access.log
stderr_logfile=/var/log/supervisor/%(program_name)s-error.log
autostart=true
autorestart=true
stopasgroup=true
