# config/supervisor/conf.d/uniboard-test.conf
[program:uniboard-test]
command=/opt/uniboard/bin/python /opt/uniboard/src/uniboard/portality/app.py
user=cloo
directory=/opt/uniboard/src/uniboard
stdout_logfile=/var/log/supervisor/%(program_name)s-access.log
stderr_logfile=/var/log/supervisor/%(program_name)s-error.log
autostart=true
autorestart=true
stopasgroup=true
